# docs/sensor-kickstart/com.studnicky.yamete.sensor-kickstart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.studnicky.yamete.sensor-kickstart</string>

    <!--
        Long-lived daemon mode. The helper runs `kickstart()` once on
        startup and then subscribes to IOKit system power notifications
        via `IORegisterForSystemPower`, re-running the kickstart on every
        `kIOMessageSystemHasPoweredOn` event so the sensor stays active
        across sleep / wake cycles even on hardware or macOS revisions
        where the driver cools the sensor during sleep. Idle CPU cost: ~0
        (it is parked in `CFRunLoopRun`).
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/libexec/yamete-sensor-kickstart</string>
        <string>daemon</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!--
        Restart if the daemon ever exits (it shouldn't — `CFRunLoopRun`
        blocks indefinitely — but if something crashes it, come back up).
    -->
    <key>KeepAlive</key>
    <true/>

    <!--
        If the daemon crashes in a tight loop, throttle the restart rate
        so we don't spin launchd. 10 seconds between respawns.
    -->
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <!--
        Background process type: lower scheduler / I/O priority than
        Interactive or Standard. Appropriate for a passive watcher.
    -->
    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/var/log/yamete-sensor-kickstart.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/yamete-sensor-kickstart.log</string>
</dict>
</plist>
